# utils/supervisor/k8s/management.conf
# Copyright (c) ASTRID 2020-2022
# author: Alex Carrega <alessandro.carrega@cnit.it>

[eventlistener:multivisor-rpc]
command=multivisor-rpc --bind 0:9002
events=PROCESS_STATE,SUPERVISOR_STATE_CHANGE
autostart=true
autorestart=true
stdout_logfile=/var/log/management/multivisor-rpc.out.log
stderr_logfile=/var/log/management/multivisor-rpc.err.log

[eventlistener:supervisor-alert]
command=supervisor-alert --telegram --show-hostname
events=PROCESS_STATE_RUNNING,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
autostart=true
autorestart=true
stdout_logfile=/var/log/management/supervisor-alert.out.log
stderr_logfile=/var/log/management/supervisor-alert.err.log

[program:multivisor]
directory=/root/guard-platform/utils/supervisor
command=multivisor -c multivisor.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/management/multivisor.out.log
stderr_logfile=/var/log/management/multivisor.err.log

[ctlplugin:wildcards]
supervisor.ctl_factory = supervisorwildcards.controllerplugin:make_wildcards_controllerplugin
match_group = 1

[eventlistener:mem-mon@cb_manager-port_forward]
command=memmon -p cb_manager-port_forward=200MB -m alessandro.carrega@cnit.it
events=TICK_60
autostart=true
autorestart=true
stdout_logfile=/var/log/management/mem-mom@cb_manager-port_forward.out.log
stderr_logfile=/var/log/management/mem-mom@cb_manager-port_forward.err.log

[eventlistener:mem-mon@elasticsearch-port_forward]
command=memmon -p elasticsearch-port_forward=200MB -m alessandro.carrega@cnit.it
events=TICK_60
autostart=true
autorestart=true
stdout_logfile=/var/log/management/mem-mom@elasticsearch-port_forward.out.log
stderr_logfile=/var/log/management/mem-mom@elasticsearch-port_forward.err.log

[eventlistener:fatal-mail-batch]
command=fatalmailbatch --toEmail="alessandro.carrega@cnit.it" --fromEmail="astrid@k8s.cnit" --subject="Crash report from astrid@k8s.cnit"
events=PROCESS_STATE,TICK_60
autostart=true
autorestart=true
stdout_logfile=/var/log/management/fatal-mail-batch.out.log
stderr_logfile=/var/log/management/fatal-mail-batch.err.log

[eventlistener:crash-mail-batch]
command=crashmailbatch --toEmail="alessandro.carrega@cnit.it" --fromEmail="astrid@k8s.cnit" --subject="Crash report from astrid@k8s.cnit"
events=PROCESS_STATE,TICK_60
autostart=true
autorestart=true
stdout_logfile=/var/log/management/crash-mail-batch.out.log
stderr_logfile=/var/log/management/crash-mail-batch.err.log

[ctlplugin:quick]
supervisor.ctl_factory = supervisor_quick:make_quick_controllerplugin
